# config/.crontab
*/1 * * * * /usr/bin/node /var/www/MuLetter/muletter-be/dist/utils/errorRetry >> /var/www/MuLetter/muletter-be/logs/retry/log_`date +\%Y\%m\%d`_`date +\%H\%M\%S`.log 2>&1